# contrib/init/baricoind.conf
description "Baricoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env BARICOIND_BIN="/usr/bin/baricoind"
env BARICOIND_USER="baricoin"
env BARICOIND_GROUP="baricoin"
env BARICOIND_PIDDIR="/var/run/baricoind"
# upstart can't handle variables constructed with other variables
env BARICOIND_PIDFILE="/var/run/baricoind/baricoind.pid"
env BARICOIND_CONFIGFILE="/etc/baricoin/baricoin.conf"
env BARICOIND_DATADIR="/var/lib/baricoind"

expect fork

respawn
respawn limit 5 120
kill timeout 600

pre-start script
    # this will catch non-existent config files
    # baricoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$BARICOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run baricoind."
        echo "The setting must appear in $BARICOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Baricoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$BARICOIND_PIDDIR"
    chmod 0755 "$BARICOIND_PIDDIR"
    chown $BARICOIND_USER:$BARICOIND_GROUP "$BARICOIND_PIDDIR"
    chown $BARICOIND_USER:$BARICOIND_GROUP "$BARICOIND_CONFIGFILE"
    chmod 0660 "$BARICOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$BARICOIND_PIDFILE" \
    --chuid $BARICOIND_USER:$BARICOIND_GROUP \
    --exec "$BARICOIND_BIN" \
    -- \
    -pid="$BARICOIND_PIDFILE" \
    -conf="$BARICOIND_CONFIGFILE" \
    -datadir="$BARICOIND_DATADIR" \
    -disablewallet \
    -daemon

